# jandroid.service
[Unit]
Description=Jandroid

[Service]
ExecStart=/home/jan/JandroidC/jandroid

[Install]
WantedBy=multi-user.target
